# target/linux/mediatek/filogic/base-files/etc/init.d/smp_affinity
#!/bin/sh /etc/rc.common

START=24

start() {
	local board=$(board_name)
	local cpu=` grep -c ^processor /proc/cpuinfo `
	local ethdevice='15100000.ethernet'
	local usbdevice='xhci-hcd:usb1'
	local wifidevice='mt7915e'
	
	case $board in
	*gl-mt6000)
        IRQs=`ls -d /proc/irq/*/$ethdevice | awk -F\/ '{ print $4 }'`
			for ETHirq in ${IRQs}; do
				echo $(( 2**(${count}+1))) > /proc/irq/${ETHirq}/smp_affinity
				let "count=${count}+2"
				let "cpu=${count}-1"
				logger -s "smp_affinity: assigning ${ethdevice} to CPU${cpu} on IRQ ${ETHirq}"
			done
		USBirq=`ls -d /proc/irq/*/$usbdevice | awk -F\/ '{ print $4 }'`
		USBcpu=`cat /proc/irq/${USBirq}/smp_affinity`
		if [ "$USBcpu" == "f" ]; then
			logger -s "smp_affinity: leave ${usbdevice} to All CPUs on IRQ ${USBirq}"
		else
			logger -s "smp_affinity: leave ${usbdevice} to CPU${USBcpu} on IRQ ${USBirq}"
		fi
		WIFIirq=`ls -d /proc/irq/*/$wifidevice | awk -F\/ '{ print $4 }'`
		WIFIcpu=`cat /proc/irq/${WIFIirq}/smp_affinity`
		if [ "$WIFIcpu" == "f" ]; then
			logger -s "smp_affinity: leave ${wifidevice} to All CPUs on IRQ ${WIFIirq}"
		else
			logger -s "smp_affinity: leave ${wifidevice} to CPU${WIFIcpu} on IRQ ${WIFIirq}"
		fi
    ;;
    *)
		logger -s "smp_affinity:  Not compatible with ${BOARD}"
	;;
    esac
	
}

stop() {
	local board=$(board_name)
	local ethdevice='15100000.ethernet'
	
	case $board in
	*gl-mt6000)
        IRQs=`ls -d /proc/irq/*/$ethdevice | awk -F\/ '{ print $4 }'`
			for ETHirq in ${IRQs}; do
				echo f > /proc/irq/${ETHirq}/smp_affinity
				logger -s "smp_affinity: restore ${ethdevice} to AllCPUs on IRQ ${ETHirq}"
			done
    ;;
    *)
		logger -s "smp_affinity:  Not compatible with ${BOARD}"
	;;
    esac
	
}
